# todo-notifyd.service
# ~/.config/systemd/user/todo-daemon.service
[Unit]
Description=TODO Notification Daemon
After=graphical-session.target

[Service]
Type=simple
# Adjust these paths as needed
ExecStart=/usr/bin/python3 %h/.local/bin/todo-daemon.py %h/Documents %h/notes
Restart=always
RestartSec=60
Environment=DISPLAY=%i
Environment=XDG_RUNTIME_DIR=%t

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
